# noteforai.service
[Unit]
Description=NoteForAI
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/noteforai
ExecStart=/opt/noteforai/noteforai serve
Environment=PORT=8080
Environment=DATA_DIR=/opt/noteforai/data
Environment=QUOTA_MB=500
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
